# overlay/usr/local/lib/systemd/system/jackd.target
[Unit]
Description=jackd sound server target
After=jackd.service 
BindsTo=jackd.service
