# overlay/usr/local/lib/systemd/system/mn_core_state.timer
[Unit]
Description=[mn] mn_get_throttled timer

[Timer]
OnCalendar=

[Install]
WantedBy=multi-user.target
