# overlay/usr/local/lib/systemd/system/mn_get_throttled.timer
[Unit]
Description=[mn] mn_get_throttled timer

[Timer]
OnCalendar=*-*-* *:0,10,20,30,40,50:00

[Install]
WantedBy=multi-user.target
