# overlay/usr/local/lib/systemd/system/mn_jackd_watch.path
[Unit]
Description=[mn] jackd socket watcher

[Path]
PathExistsGlob=/dev/shm/jack_*_????_?

[Install]
WantedBy=multi-user.target
