# overlay/usr/local/lib/systemd/system/mn_jackd.target
[Unit]
Description=jackd sound server target
After=mn_jackd.service 
BindsTo=mn_jackd.service
