# overlay/usr/lib/systemd/system/jackd.target
[Unit]
Description=jackd sound server target
After=sound.target jackd.service 
Requires=sound.target jackd.service
