# systemd/ntp-synced.target
[Unit]
Description=System time is monotonic after ntp synchronization
After=ntp-sync.service
